# .config/autostart/01-kdeconnectd.desktop
[Desktop Entry]
Name=KDE Connect Daemon
Comment=KDE Connect Daemon
Exec=/usr/lib/kdeconnectd
Terminal=false
Type=Application
StartupNotify=true
