# .config/autostart/90-blueman-applet.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/blueman-applet
Name=blueman-applet
Comment=Blueman applet for the system tray
