# .config/autostart/90-feh.desktop
[Desktop Entry]
Name=feh
Exec=/usr/bin/feh --bg-scale .config/i3/drum_wallpaper1.jpg
Type=Application
Comment=image viewer and cataloguer
