# .config/autostart/90-feh.desktop
[Desktop Entry]
Name=feh
Exec=/bin/bash /home/danilo/tools/utils/wallpaper.sh
Type=Application
Comment=image viewer and cataloguer
