# .config/autostart/90-nm-applet.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/nm-applet
Name=nm-applet
Comment=NetworkManager applet for the system tray
