# .config/autostart/90-pasystray.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/pasystray
Name=pasystray
Comment=PulseAudio applet for the system tray
