# .config/autostart/99-nm-applet.desktop
[Desktop Entry]
Type=Application
Comment=network monitor and control GUI applet
Name=nm-applet
Exec=/usr/bin/nm-applet
