# .config/autostart/99-pasystray.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/pasystray
Name=pasystray
Comment=PulseAudio controller for the system tray
