# .config/autostart/99-ssh_agent.desktop
[Desktop Entry]
Type=Application
Exec=/home/danilo/tools/utils/launch_ssh_agent.sh
Name=SSH Agent
